# supervisord.conf
[unix_http_server]
file = /tmp/supervisor5.sock   ; (the path to the socket file)

[inet_http_server]
port = *:49001        ; (ip_address:port specifier, *:port for all iface)

[supervisord]
logfile = /tmp/supervisord5.log ; (main log file;default /supervisord.log)
logfile_maxbytes = 50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups = 10           ; (num of main logfile rotation backups;default 10)
loglevel = info                ; (log level;default info; others: debug,warn,trace)
pidfile = /tmp/supervisord5.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon = false               ; (start in foreground if true;default false)
minfds = 1024                  ; (min. avail startup file descriptors;default 1024)
minprocs = 200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = http://0.0.0.0:49001 ; use an http:// url to specify an inet socket

[program:QA41_v1_5_1_qa_train_49101]
command = gunicorn -c conf/train_gunicorn_conf.py -b 0.0.0.0:49101 train_server:app
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/train.log
stderr_logfile = log/train.log

[program:QA41_v1_5_1_qa_predict_knowledge_49102]
command = gunicorn -c conf/predict_knowledge_gunicorn_conf.py -b 0.0.0.0:49102 "predict_server:get_app('knowledge', '../core/predict/config.ini', service_id=41)"
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/predict_knowledge.log
stderr_logfile = log/predict_knowledge.log
port = 49102
since = 2019/11/11
until = 2021/12/12
ID = 41

[program:QA42_v1_5_1_qa_train_49103]
command = gunicorn -c conf/train_gunicorn_conf.py -b 0.0.0.0:49103 train_server:app
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/train.log
stderr_logfile = log/train.log

[program:QA42_v1_5_1_qa_predict_knowledge_49104]
command = gunicorn -c conf/predict_knowledge_gunicorn_conf.py -b 0.0.0.0:49104 "predict_server:get_app('knowledge', '../core/predict/config.ini', service_id=42)"
environment = FEAT_ID=42,CUDA_VISIBLE_DEVICES='2'
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/predict_knowledge.log
stderr_logfile = log/predict_knowledge.log

[program:QA43_v1_5_1_qa_train_49105]
command = gunicorn -c conf/train_gunicorn_conf.py -b 0.0.0.0:49105 train_server:app
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/train.log
stderr_logfile = log/train.log

[program:QA43_v1_5_1_qa_predict_knowledge_49106]
command = gunicorn -c conf/predict_knowledge_gunicorn_conf.py -b 0.0.0.0:49106 "predict_server:get_app('knowledge', '../core/predict/config.ini', service_id=43)"
environment = FEAT_ID=43,CUDA_VISIBLE_DEVICES='2'
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/predict_knowledge.log
stderr_logfile = log/predict_knowledge.log

[program:QA44_v1_5_1_qa_train_49107]
command = gunicorn -c conf/train_gunicorn_conf.py -b 0.0.0.0:49107 train_server:app
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/train.log
stderr_logfile = log/train.log

[program:QA44_v1_5_1_qa_predict_knowledge_49108]
command = gunicorn -c conf/predict_knowledge_gunicorn_conf.py -b 0.0.0.0:49108 "predict_server:get_app('knowledge', '../core/predict/config.ini', service_id=44)"
environment = FEAT_ID=44,CUDA_VISIBLE_DEVICES='2'
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/predict_knowledge.log
stderr_logfile = log/predict_knowledge.log

[program:nlu_train_v2_1_1_49201]
command = bash -c "bin/lu-server -J-Xmx8G -Dconfig.file=conf/prod.conf -Dpidfile.path=/dev/null -Dhttp.port=49201"
environment = LD_LIBRARY_PATH=./license
directory = ../nlu/train_server
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/nlu_train.log
stderr_logfile = log/nlu_train.log

[program:nlu_predict_v2_1_1_49202]
command = bash -c "bin/lu-server -J-Xmx8G -Dconfig.file=conf/prod.conf -Dpidfile.path=/dev/null -Dhttp.port=49202"
environment = LD_LIBRARY_PATH=./license
directory = ../nlu/predict_server
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/nlu_predict.log
stderr_logfile = log/nlu_predict.log

[program:nlu_train_v2_1_1_49203]
command = bash -c "bin/lu-server -J-Xmx8G -Dconfig.file=conf/prod.conf -Dpidfile.path=/dev/null -Dhttp.port=49203"
environment = LD_LIBRARY_PATH=./license
directory = ../nlu/train_server
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/nlu_train.log
stderr_logfile = log/nlu_train.log

[program:nlu_predict_v2_1_1_49204]
command = bash -c "bin/lu-server -J-Xmx8G -Dconfig.file=conf/prod.conf -Dpidfile.path=/dev/null -Dhttp.port=49204"
environment = LD_LIBRARY_PATH=./license
directory = ../nlu/predict_server
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/nlu_predict.log
stderr_logfile = log/nlu_predict.log

[program:ec_predict_v2_1_1_49301]
command = bash -c "bin/spellcheck-rest -Dhttp.port=49301 -Dconfig.file=conf/prod.conf -Dpidfile.path=/dev/null -J-Xms2G -J-Xmx4G"
environment = LD_LIBRARY_PATH=./data/license/
directory = ../error-correction/SpellCheck-REST/
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/ec_predict.log
stderr_logfile = log/ec_predict.log

[program:nlu_train_log]
command = bash -c "./watch_nlu_log.sh train_server"
--tries = 3
directory = .
autostart = false
autorestart = false
user = ecp
stopasgroup = true
killasgroup = true

[program:nlu_predict_log]
command = bash -c "./watch_nlu_log.sh predict_server"
--tries = 3
directory = .
autostart = false
autorestart = false
user = ecp
stopasgroup = true
killasgroup = true

[program:nvidia_watching]
command = bash -c "nvidia-smi dmon -i 1 -s mu -d 5 -o TD"
autostart = false
autorestart = false
user = ecp
stopasgroup = true
killasgroup = true

[program:ram_disk_watching]
command = bash -c "./watch_df.sh"
directory = .
autostart = false
autorestart = false
user = ecp
stopasgroup = true
killasgroup = true

[program:QA45_v1_5_1_qa_train_49109]
command = gunicorn -c conf/train_gunicorn_conf.py -b 0.0.0.0:49109 train_server:app
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/train.log
stderr_logfile = log/train.log

[program:QA45_v1_5_1_qa_predict_knowledge_49110]
command = gunicorn -c conf/predict_knowledge_gunicorn_conf.py -b 0.0.0.0:49110 "predict_server:get_app('knowledge', '../core/predict/config.ini', service_id=45)"
environment = FEAT_ID=45,CUDA_VISIBLE_DEVICES='2'
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/predict_knowledge.log
stderr_logfile = log/predict_knowledge.log

[program:nlu_train_v2_1_1_49213]
command = bash -c "bin/lu-server -J-Xmx8G -Dconfig.file=conf/prod.conf -Dpidfile.path=/dev/null -Dhttp.port=49213"
environment = LD_LIBRARY_PATH=./license
directory = ../nlu/train_server
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/nlu_train.log
stderr_logfile = log/nlu_train.log

[program:nlu_predict_v2_1_1_49214]
command = bash -c "bin/lu-server -J-Xmx8G -Dconfig.file=conf/prod.conf -Dpidfile.path=/dev/null -Dhttp.port=49214"
environment = LD_LIBRARY_PATH=./license
directory = ../nlu/predict_server
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/nlu_predict.log
stderr_logfile = log/nlu_predict.log

[program:QA46_v1_5_1_qa_train_49111]
command = gunicorn -c conf/train_gunicorn_conf.py -b 0.0.0.0:49111 train_server:app
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/train.log
stderr_logfile = log/train.log

[program:QA46_v1_5_1_qa_predict_knowledge_49112]
command = gunicorn -c conf/predict_knowledge_gunicorn_conf.py -b 0.0.0.0:49112 "predict_server:get_app('knowledge', '../core/predict/config.ini', service_id=46)"
environment = FEAT_ID=46
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/predict_knowledge.log
stderr_logfile = log/predict_knowledge.log

[program:QA47_v1_5_1_qa_train_49113]
command = gunicorn -c conf/train_gunicorn_conf.py -b 0.0.0.0:49113 train_server:app
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/train.log
stderr_logfile = log/train.log

[program:QA47_v1_5_1_qa_predict_knowledge_49114]
command = gunicorn -c conf/predict_knowledge_gunicorn_conf.py -b 0.0.0.0:49114 "predict_server:get_app('knowledge', '../core/predict/config.ini', service_id=47)"
environment = FEAT_ID=47
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/predict_knowledge.log
stderr_logfile = log/predict_knowledge.log

[program:QA48_v1_5_1_qa_train_49115]
command = gunicorn -c conf/train_gunicorn_conf.py -b 0.0.0.0:49115 train_server:app
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/train.log
stderr_logfile = log/train.log

[program:QA48_v1_5_1_qa_predict_knowledge_49116]
command = gunicorn -c conf/predict_knowledge_gunicorn_conf.py -b 0.0.0.0:49116 "predict_server:get_app('knowledge', '../core/predict/config.ini', service_id=48)"
environment = FEAT_ID=48
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/predict_knowledge.log
stderr_logfile = log/predict_knowledge.log

[program:testt]
command = gunicorn -c conf/predict_knowledge_gunicorn_conf.py -b 0.0.0.0:9111 "predict_server:get_app('knowledge', '../core/predict/config.ini', service_id=1)"
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/predict_knowledge.log
stderr_logfile = log/predict_knowledge.log
port = 9111
since = 2019/01/01
until = 2021/12/12
ID = 1

[program:test1111]
command = gunicorn -c conf/train_gunicorn_conf.py -b 0.0.0.0:9100 train_server:app
directory = .
startsecs = 1
stopwaitsecs = 1
autostart = false
autorestart = false
stdout_logfile = log/train.log
stderr_logfile = log/train.log
since = 2019/01/10
until = 2020/12/12
port = 9100
ID = 2

